# lib/mac/crontab
# DFK crontab
#=============================================
#  Each crontab file entry consists of a line with six fields, separated by
#  spaces and tabs.  The information in the fields specifies when the command
#  runs and the name of the command.  The fields specify the following:
#
#    o  The first field specifies the minute (0 to 59).
#
#    o  The second field specifies the hour (0 to 23).
#
#    o  The third field specifies the day of the month (1 to 31).
#
#    o  The fourth field specifies the month of the year (1 to 12).
#
#    o  The fifth field specifies the day of the week (0 to 6 for Sunday to
#       Saturday).
#
#    o  The sixth field specifies the shell command to be executed.
#
#  You can specify the following values in the fields that indicate the time:
#
#    o  An integer (within the appropriate range of values)
#
#    o  Two integers separated by a dash to indicate an inclusive range
#
#    o  A list of integers separated by commas
#
#    o  An asterisk to select all possible values
#=============================================

PATH=/usr/local/bin:/usr/bin:/bin:/u/dfk/scripts/common:/u/dfk/scripts/mac
# MAILTO=David.F.Kotz@Dartmouth.edu

# don't forget my crontab
0 7 * * * cd /u/dfk/lib/mac && crontab -l > crontab && git-commit-quietly "crontab" crontab

# don't forget my hosts file
4 7 * * * cd /u/dfk/lib/mac && cat /etc/hosts > etc-hosts && git-commit-quietly "crontab" etc-hosts

# track what I have installed with homebrew
8 7 * * * cd /u/dfk/lib/mac && brew ls > brew-list && git-commit-quietly "crontab" brew-list

# pull any changes to lib, dot, scripts
0 9 * * *   cd projects/kotz-bin-lib && git pull --quiet

# show me what files are left uncommitted
2 9 * * *   git-to-commit projects/* 
4 9 * * *   svn-to-commit projects/* 
# 6 9 * * * cvs-to-commit projects/*

